# deployment/open511-importer.conf
description "open511-importer"

start on runlevel [2345]
stop on runlevel [016]

respawn

script
	cd /home/open511/mtl511
	su -c "/home/open511/env-geotrafic/bin/python import_runner.py" open511
end script
