# deployment/open511-importer.conf
description "open511-importer"

start on runlevel [2345]
stop on runlevel [016]

respawn

script
	cd /home/opendev/mtl511
	su -c "/home/opendev/env-geotrafic/bin/python import_runner.py" opendev
end script
