# deployment/open511-web.conf
description "open511-web-server"

start on runlevel [2345]
stop on runlevel [016]

respawn

script
	cd /home/opendev/mtl511
	su -c "/home/opendev/env-geotrafic/bin/gunicorn -c gunicorn_settings.py geotrafic511.wsgi" opendev
end script
