# deployment/open511-web.conf
description "open511-web-server"

start on (filesystem)
stop on runlevel [016]

respawn
setuid opendev
chdir /home/opendev/mtl511

exec /home/opendev/env-geotrafic/bin/gunicorn -c gunicorn_settings.py geotrafic511.wsgi
